# etc/udev/rules.d/30-leds.rules
ACTION=="add", RUN+="/bin/sh -c 'chown droidian:droidian /sys/class/leds/*/*'"
